# roles/nwn2server/templates/backup@.service.j2
[Unit]
Description=Backup NWN2 servervault, bioware database and MySQL db

[Service]
Type=oneshot
ExecStart=/srv/%i/bin/backup